# container-files/etc/services.d/zabbix-proxy/run
#!/bin/execlineb -P

s6-setuidgid zabbix
/usr/sbin/zabbix_proxy --config /etc/zabbix/zabbix_proxy.conf
